# config/systemd/user/battery-notify.service
[Unit]
Description=Battery Notifier

[Service]
Type=oneshot
ExecStart=/bin/bash -c "/home/ajdiaz/.config/i3/battery-notify"
