# config/systemd/user/batt_notify.service
# offlineimap.service

[Unit]
Description=Battery Notifier

[Service]
Type=oneshot
ExecStart=/home/ajdiaz/.bin/batt_notify
