# config/systemd/user/batt_notify.timer
# offlineimap.timer

[Unit]
Description=Battery Notifier Timer

[Timer]
OnCalendar=*:0/2
Unit=batt_notify.service

[Install]
WantedBy=default.target
