# config/systemd/user/system-status.timer
[Unit]
Description=System Status

[Timer]
OnCalendar=hourly
Unit=system-status.service
RandomizedDelaySec=60

[Install]
WantedBy=default.target
