# config/systemd/user/system-status.timer
[Unit]
Description=System Status

[Timer]
OnCalendar=*:0/60
Unit=system-status.service

[Install]
WantedBy=default.target
